# debian/keyboard-configuration.preinst
#!/bin/sh

set -e

for file in keyboard-setup console-setup; do
    dpkg-maintscript-helper rm_conffile /etc/init.d/$file 1.138~ -- "$@"
    if [ -x /etc/init.d/$file ]; then
        update-rc.d $file remove >/dev/null
    fi
done

#DEBHELPER#
